# contrib/upstart/cjdns-dynamic.conf
description "dynamic DNS peer resolver for cjdns"

author "Adam Novak <interfect@gmail.com>"
# Based on cjsdns Upstart config by Sergey Davidoff <shnatsel@gmail.com>
# Requires that you be using Upstart to control cjdns

start on started cjdns
stop on stopping cjdns

# respawn if it crashes, but no more than once in 180 seconds
respawn
respawn limit 2 300

pre-start exec /usr/lib/cjdns/dynamicEndpointsGenConf.sh

exec /usr/lib/cjdns/dynamicEndpoints.py /etc/cjdns-dynamic.conf --adminInfo /etc/cjdns-dynamic.cjdnsadmin
